# xgboost-daily.service
[Unit]
Description=XGBoost Daily Trading Pipeline
After=network.target mysql.service
Wants=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/www/wwwroot/xgboost-qcv2
Environment=PYTHONPATH=/www/wwwroot/xgboost-qcv2
ExecStart=/usr/bin/python3 /www/wwwroot/xgboost-qcv2/multi_day_runner.py --days 1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xgboost-daily

# Resource limits
LimitNOFILE=65536
MemoryMax=4G
CPUQuota=200%

[Install]
WantedBy=multi-user.target